# EdgeEye_Client_v1.0.1/EdgeEye_1.0.1/Standard/bin/iCAP_ClientService/linux/bin/iCAP_ClientService_linux/iCAP_Client.service
[Unit]
Description=iCAP Client Service

[Service]
Type=simple
Restart=always
ExecStart=/usr/sbin/iCAP_ClientService

[Install]
WantedBy=multi-user.target
